# src/Docker/Services/Worker/fixtures/supervisord.conf
[supervisord]
nodaemon=true

[program:worker]
process_name=%(program_name)s
command=php /var/www/artisan queue:work
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/storage/logs/queue_worker.log
